# scripts/launchd/com.amclinica.local3000.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.amclinica.local3000</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-lc</string>
        <string>/Users/arimacm5/.am-clinica/start-local-3000.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>30</integer>
    <key>StandardOutPath</key>
    <string>/Users/arimacm5/Library/Logs/am-clinica-local3000.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/arimacm5/Library/Logs/am-clinica-local3000.error.log</string>
</dict>
</plist>
